# scripts/etc/cron.hourly/shark-statsconv-update_ip
#!/bin/sh

# Update IP conversation stats under /cap/stats

/usr/local/bin/shark-statsconv ip >> /var/log/shark-statsconv.log

